# host/docker.service
[Service]
ExecStart=
ExecStart=/usr/bin/docker -d -r=false --icc=false
