# arch-nginx/supervisord/nginx.ini
[program:nginx]
command=/usr/bin/nginx
priority=100
startsecs=5
stopsignal=QUIT
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/%(program_name)s.stdout.log
stderr_logfile=/var/log/supervisor/%(program_name)s.stderr.log
